# deploy/elfie-agent.service
[Unit]
# Template — install_services.sh fills in the __ELFIE_*__ placeholders for
# your machine. Don't edit this copy directly; run the installer.
Description=Elfie voice agent worker
After=network-online.target docker.service

[Service]
Environment=PATH=__ELFIE_PATH__
WorkingDirectory=__ELFIE_DIR__
ExecStart=__ELFIE_DIR__/.venv/bin/python elfie_agent.py dev
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
